# ansible/roles/app_servers/templates/app.service.j2
[Unit]
Description=Flask Web Application
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/app
Environment=PATH=/opt/app/venv/bin
Environment=MYSQL_HOST={{ db_host }}
Environment=MYSQL_DATABASE={{ db_name }}
Environment=MYSQL_USER={{ db_user }}
Environment=MYSQL_PASSWORD={{ db_password }}
ExecStart=/opt/app/venv/bin/python app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target